# docs/xcaox-maps.service
# XCAOX Maps Systemd Service
# Place this file at: /etc/systemd/system/xcaox-maps.service
# Alternative to PM2 for system-level service management

[Unit]
Description=XCAOX Maps - Interactive Mapping Platform
Documentation=https://github.com/your-username/xcaox-maps
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/xcaox-maps
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
SyslogIdentifier=xcaox-maps

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/xcaox-maps
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target

# Usage:
# sudo systemctl enable xcaox-maps
# sudo systemctl start xcaox-maps
# sudo systemctl status xcaox-maps
# sudo journalctl -u xcaox-maps -f
